# files/etc/supervisor.d/hg.ini
[program:hg]
command=uwsgi --plugins http,python --http-socket :9002 --wsgi-file /srv/hg/uwsgi.py --uid www --gid www --enable-threads --thunder-lock --pidfile /var/run/hg.pid --master --buffer-size=32768
priority=100
autostart=true
autorestart=true
stdout_logfile=/proc/1/fd/1
stderr_logfile=/proc/1/fd/2
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0